# etc/xdg/autostart/lxpanel-restart.desktop
[Desktop Entry]
Name=LxpanelRestart
Categories=Utility;
Type=Application
Exec=lxpanel-restart
Terminal=false
NoDisplay=false
